# ws/ws-proxy.service
[Unit]
Description=WebSocket to TCP Proxy Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/usr/local/bin
ExecStart=/usr/bin/node /usr/local/bin/proxy.js
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
